# buildroot-overlay/etc/profile.d/webos-env.sh
export WEBOS_HOME=/opt/webos
export WEBOS_DATA=/opt/webos/data
