# live-build/config/includes.chroot/etc/systemd/system/myos.service
[Unit]
Description=Start MyOS desktop session
After=systemd-user-sessions.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=myos
Group=myos
Environment=HOME=/home/myos
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/opt/myos/start.sh
Restart=on-failure
RestartSec=2
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
TTYVTDisallocate=yes
StandardInput=tty
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
